# supervisord.conf
[program:logstash-forwarder]
command=/opt/logstash-forwarder/logstash-forwarder -config /etc/logstash/config.json
numprocs=1
autostart=true
autorestart=true
